# mcp-servers/automation-mcp/automation-mcp.service
[Unit]
Description=Workflow automation MCP
After=network.target

[Service]
Type=simple
User=matt-woodworth
WorkingDirectory=/home/matt-woodworth/dev/myroofgenius-backend/mcp-servers/automation-mcp
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/matt-woodworth/dev/myroofgenius-backend/mcp-servers/automation-mcp/server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
