# mcp-servers/database-mcp/database-mcp.service
[Unit]
Description=Database operations MCP server
After=network.target

[Service]
Type=simple
User=mwwoodworth
WorkingDirectory=/home/mwwoodworth/code/mcp-servers/database-mcp
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/mwwoodworth/code/mcp-servers/database-mcp/server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
