# mcp-servers/mcp-servers.service
#!/bin/bash
# MCP Servers Systemd Service Script
# This script manages all MCP servers as a single service

[Unit]
Description=BrainOps MCP Servers
After=network.target

[Service]
Type=forking
User=mwwoodworth
WorkingDirectory=/home/mwwoodworth/code/mcp-servers
ExecStart=/home/mwwoodworth/code/mcp-servers/start_all_mcp.sh
ExecStop=/home/mwwoodworth/code/mcp-servers/stop_all_mcp.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target