# mcp-servers/mcp-servers.service
[Unit]
Description=BrainOps MCP Servers
After=network.target

[Service]
Type=forking
User=matt-woodworth
WorkingDirectory=/home/matt-woodworth/dev/myroofgenius-backend/mcp-servers
ExecStart=/home/matt-woodworth/dev/myroofgenius-backend/mcp-servers/start_all_mcp.sh
ExecStop=/home/matt-woodworth/dev/myroofgenius-backend/mcp-servers/stop_all_mcp.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target