# mcp-servers/ai-orchestrator-mcp/ai-orchestrator-mcp.service
[Unit]
Description=AI agent orchestration MCP
After=network.target

[Service]
Type=simple
User=matt-woodworth
WorkingDirectory=/home/matt-woodworth/dev/myroofgenius-backend/mcp-servers/ai-orchestrator-mcp
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/matt-woodworth/dev/myroofgenius-backend/mcp-servers/ai-orchestrator-mcp/server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
